# public_html/python/investment-mcp.service
[Unit]
Description=Investment Portfolio MCP Server
After=network.target mysql.service mariadb.service
Wants=mysql.service mariadb.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/html/investments.davidappleyard.net

# Load credentials from .env — adjust path if your .env lives elsewhere
EnvironmentFile=/var/www/html/investments.davidappleyard.net/.env

ExecStart=/opt/investment-mcp-venv/bin/python3 /var/www/html/investments.davidappleyard.net/public_html/python/mcp_server.py
Restart=on-failure
RestartSec=10

# Keep stdout/stderr in journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=investment-mcp

[Install]
WantedBy=multi-user.target
